# zabbix-proxy.service
[Unit]
Description=Zabbix Proxy
After=network.target

[Service]
Type=simple
ExecStart=/usr/sbin/zabbix_proxy -f
User=root

[Install]
WantedBy=multi-user.target
